# vexpress_tools/rootfs/etc/profile.d/alias.sh
#!/bin/sh
# Common command aliases

alias ll='ls -alF --color=auto'
alias la='ls -A'
alias l='ls -CF'
alias grep='grep --color=auto'
