# vexpress_tools/rootfs/etc/profile.d/banner.sh
#!/bin/sh
echo "Welcome $(whoami)! The system uptime is: $(uptime -p)"
