# vexpress_tools/rootfs/etc/inittab
::sysinit:/etc/init.d/rcS
# For QEMU serial
#::respawn:/sbin/getty 115200 console
::askfirst:/bin/sh
